# src/qz/installer/assets/mac-launchagent.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key><string>%PACKAGE_NAME%</string>
        <key>KeepAlive</key>
        <dict>
            <key>SuccessfulExit</key><false/>
            <key>AfterInitialDemand</key><false/>
        </dict>
        <key>RunAtLoad</key><true/>
        <key>ProgramArguments</key>
        <array>
            <string>%COMMAND%</string>
            <string>%PARAM%</string>
        </array>
    </dict>
</plist>